# systemd/dropins/livekit-server.service.d/wait-tailnet.conf
# Gate startup until the tailscale0 100.x address exists (livekit.yaml binds it).
# Waits up to 60s, then fails fast so systemd surfaces the missing tailnet.
[Service]
ExecStartPre=/bin/sh -c 'for i in $(seq 1 60); do ip -4 addr show tailscale0 2>/dev/null | grep -q "100\\." && exit 0; sleep 1; done; echo "tailscale0 100.x address never appeared" >&2; exit 1'
