# systemd/dropins/skchat-nostr-relay.service.d/override.conf
# Host-specific: bind the relay to this host's tailnet IP so tailnet peers can
# reach it. The value below is the .158 tailscale0 address; EDIT IT for a new
# host (or template it from `tailscale ip -4`). Not a secret.
[Service]
Environment=NOSTR_RELAY_HOST=100.108.59.57
