# systemd/skchat-lumina-call.service.d/musetalk.conf
[Service]
Environment=LUMINA_MUSETALK_URL=http://chiap08:18803
Environment=LUMINA_MUSETALK_REFERENCE=/home/skuser01/sovereign-facetime/assets/lumina.png
Environment=LUMINA_MUSETALK_FPS=25
