# systemd/skchat-bridges.target
[Unit]
Description=SKChat bridge services (Lumina + Opus)
After=skchat-daemon.service
Wants=skchat-lumina-bridge.service skchat-opus-bridge.service

[Install]
WantedBy=default.target
