# systemd/dropins/skchat-webui@lumina.service.d/guest.conf
# Guest-link support for the lumina webui. Same shared secret as the lumina
# daemon guest.conf (guest links must verify across both). NOT committed;
# provision %h/.config/skchat/guest.env with:
#   SKCHAT_GUEST_TOKEN_SECRET=<secret>
# Optional (-) so a missing file does not take the webui down.
[Service]
Environment=SKCHAT_GUEST_LINKS_ENABLED=1
EnvironmentFile=-%h/.config/skchat/guest.env
