# systemd/dropins/skchat-daemon.service.d/dm-ratchet.conf
# RFC-0001 P1 --- enable the Level-3 1:1 DM ratchet (pqdr1). Capability-gated:
# only peers advertising pqdr1 are ratcheted; app/older clients stay classical.
[Service]
Environment=SKCHAT_DM_RATCHET=1
