# systemd/dropins/skchat-daemon.service.d/group.conf
[Service]
Environment=SKCHAT_GROUPS=group:c68ec8e8-f3f3-4278-af7e-fb24af6a6f4e
# DEV backend first (live skgateway model); flip SKCHAT_GROUP_MODEL to reg:ornith once ornith lands.
Environment=SKCHAT_GROUP_BACKEND_URL=http://localhost:18780/v1/chat/completions
Environment=SKCHAT_GROUP_MODEL=sk-default
Environment=SKCHAT_GROUP_ON_ERROR=note
Environment=SKCHAT_ADVOCACY_DISABLED=1
