# systemd/dropins/skchat-daemon.service.d/guest.conf
# Guest-link support for the lumina daemon. The guest token secret is NOT
# committed; it is loaded from an EnvironmentFile. Provision
# %h/.config/skchat/guest.env with:
#   SKCHAT_GUEST_TOKEN_SECRET=<secret>
# The SAME secret must be present for skchat-webui@lumina (see
# dropins/skchat-webui@lumina.service.d/guest.conf) so guest links verify across
# both. Provisioning is handled by the secrets task (P1.2). The reference is
# optional (-) so a missing file does not take the daemon down; guest links are
# simply disabled until the secret is present.
[Service]
Environment=SKCHAT_GUEST_LINKS_ENABLED=1
EnvironmentFile=-%h/.config/skchat/guest.env
