# systemd/units/skchat-telegram-lumina.service
# skchat-telegram-lumina.service: Telegram bridge for the Lumina agent.
#
# Live concrete unit reconciled from .158 (@seaBird_Lumi_bot). See the sibling
# skchat-telegram-opus.service header and README.md. Lumina's default brain is
# the abliterated qwen (her personality/creative/memory voice), pinned via
# SKC_BRIDGE_DEFAULT_ROLE=sk-creative in the override.conf drop-in. The bot
# token lives in the EnvironmentFile below (never committed).

[Unit]
Description=SKChat Telegram Bridge (Lumina agent, @seaBird_Lumi_bot)
Documentation=https://github.com/smilinTux/skchat
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=%h/.config/skchat/telegram-lumina.env
Environment=PYTHONUNBUFFERED=1
Environment=SKAGENT=lumina
Environment=SKCAPSTONE_AGENT=lumina
Environment=SKC_BRIDGE_AGENT=Lumina
Environment=SKC_BRIDGE_AGENT_HOME=%h/.skcapstone/agents/lumina
Environment=SKC_BRIDGE_LLM_URL=http://192.168.0.100:8082/v1/chat/completions
Environment=SKC_BRIDGE_LLM_MODEL=qwen3.6-27b-abliterated
Environment=SKC_BRIDGE_CTX=32768
Environment=SKC_BRIDGE_SYS_BUDGET=9000
Environment=SKC_BRIDGE_MAX_TOKENS=1024
Environment=SKCHAT_REPO=%h/clawd/skcapstone-repos/skchat
Environment=PATH=%h/.skenv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/.skenv/bin/python %h/clawd/skcapstone-repos/skchat/scripts/telegram_bridge.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=skchat-telegram-lumina

[Install]
WantedBy=default.target
