# systemd/units/skchat-app-web.service
# skchat-app-web.service: static file server for the built Flutter web client.
#
# Serves skchat-app/build/web on :8088 via scripts/serve-app-web.sh, a
# hardened stdlib-only static server (Task 7 of skchat-resilience-v1, coord
# b5078963). Replaces the earlier stopgap dev-style server
# (`python3 -m http.server`). The new server sets correct Content-Type for
# .wasm/.js/.json/.css, sends no-cache on index.html and a long immutable
# cache on content-hashed filenames, and disables directory listing
# (autoindex) entirely. Binds loopback only (127.0.0.1); the tailnet/funnel
# ingress fronts this port (see Task 4, systemd/tailscale-ingress.sh).
#
# Requires the Flutter web build to exist at SKCHAT_APP_WEB_ROOT below. Build
# it on a machine with headroom (NOT .158): see docs/deploy-plan.

[Unit]
Description=SKChat Flutter web static server (:8088)
Documentation=https://github.com/smilinTux/skchat
After=network.target

[Service]
Environment=SKCHAT_APP_WEB_ROOT=%h/clawd/skcapstone-repos/skchat-app/build/web
Environment=SKCHAT_APP_WEB_PORT=8088
Environment=SKCHAT_APP_WEB_BIND=0.0.0.0
ExecStart=/usr/bin/bash %h/clawd/skcapstone-repos/skchat/scripts/serve-app-web.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=skchat-app-web

[Install]
WantedBy=default.target
