# systemd/units/telegram-catchup.timer
# telegram-catchup.timer: fires telegram-catchup.service daily at 06:00 local
# (randomized up to 5 min). Persistent so a missed run (box asleep) fires on
# next boot. Reconciled verbatim from live .158.

[Unit]
Description=Daily Telegram Catch-Up Timer
Documentation=https://github.com/smilinTux/skchat

[Timer]
OnCalendar=*-*-* 06:00:00
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
