# systemd/dropins/skchat-lumina-call.service.d/fixes.conf
# Live fixes applied on top of the base unit. Absolute /home/cbrd21 avatar path
# replaced with %h. MuseTalk disabled (see musetalk.conf), narrate points at the
# qwen server, STT RMS floor raised.
[Service]
Environment=LUMINA_NARRATE_URL=http://192.168.0.100:8082
Environment=LUMINA_STT_MIN_RMS=350
Environment=LUMINA_MUSETALK_URL=
Environment=LUMINA_AVATAR_PATH=%h/.skcapstone/agents/lumina/avatar/portrait.png
