# systemd/dropins/skchat-lumina-call.service.d/musetalk.conf
[Service]
# MuseTalk disabled --- too slow on the current chiap08 path (~25s for typical reply).
# Re-enable by uncommenting once /ws/render streaming is wired in lumina-call.py.
Environment=LUMINA_MUSETALK_URL=
