# systemd/dropins/skchat-lumina-call.service.d/tts.conf
[Service]
Environment=SKCHAT_TTS_URL=http://localhost:15091/audio/speech
Environment=SKCHAT_TTS_VOICE=af_heart
