# systemd/dropins/skchat-lumina-call.service.d/vad.conf
[Service]
Environment=LUMINA_VAD_RMS=400
