# systemd/dropins/skchat-telegram-opus.service.d/override.conf
# Rollout flags (added 2026-07-04):
#  - inline 👍/👎 rating buttons -> ~/.skcapstone/models/ratings.jsonl -> sk-auto learns
#  - memory backend = skmem-pg (pgvector hybrid), NOT deprecated chroma/skvector
[Service]
# Lazy per-turn tool curation DISABLED (hid most tools behind keyword-matching).
# Environment=SKC_BRIDGE_LAZY_TOOLS=1
Environment=SKC_BRIDGE_RATING_BUTTONS=1
Environment=SKMEMORY_VECTOR_BACKEND=pgvector
Environment=SKMEMORY_EMBED_URL=http://192.168.0.100:11438/v1/embeddings
Environment=SKMEMORY_EMBED_MODEL=mxbai-embed-large
# SKMEMORY_PG_DSN contains a password and is NOT committed. Provision
# %h/.config/skchat/memory-pg.env with:
#   SKMEMORY_PG_DSN=postgresql://USER:PASS@HOST:5432/skmemory
# Optional (-) so a missing file does not take the bridge down; memory recall
# degrades until the DSN is present. (Rotate the inline DSN: plan task P1.2.)
EnvironmentFile=-%h/.config/skchat/memory-pg.env
