# launchd/com.skcapstone.skchat-lumina-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.skcapstone.skchat-lumina-bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>${HOME}/.skenv/bin/python3</string>
    <string>scripts/lumina-bridge.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>${HOME}/clawd/skcapstone-repos/skchat</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>LUMINA_BRIDGE_INTERVAL</key>
    <string>3</string>
    <key>SKCHAT_IDENTITY</key>
    <string>capauth:opus@skworld.io</string>
    <key>PATH</key>
    <string>${HOME}/.skenv/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>5</integer>

  <key>StandardOutPath</key>
  <string>${HOME}/.skchat/lumina-bridge.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>${HOME}/.skchat/lumina-bridge.stderr.log</string>
</dict>
</plist>
